# overlays/coreos/etc/systemd/system/scheduled-reboot.timer
[Unit]
Description=Schedules reboot

[Timer]
OnCalendar=Sat *-*-* 07:30:00

[Install]
WantedBy=timers.target
